# supporting/vnc-chrominum/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_VNC_RESOLUTION)sx%(ENV_VNC_COL_DEPTH)s
autorestart=true
user=headless
priority=100

[program:x11vnc]
command=/usr/bin/x11vnc -rfbport %(ENV_VNC_PORT)s -xkb -noxrecord -noxfixes -noxdamage -display :99 -passwd %(ENV_VNC_PW)s -wait 5 -shared
user=headless
autorestart=true
priority=200

[program:openbox]
environment=HOME="/home/headless",DISPLAY=":99",USER="headless"
command=/usr/bin/openbox
user=headless
autorestart=true
priority=300

[program:chromium]
environment=HOME="/home/headless",DISPLAY=":99",USER="headless"
command=/usr/bin/chromium-browser --no-sandbox --start-fullscreen --no-gpu "%(ENV_URL)s"
user=headless
autorestart=true
priority=400

[program:novnc]
environment=HOME="/home/headless",USER="headless"
command=/usr/bin/novnc_server --listen %(ENV_NO_VNC_PORT)s --vnc 127.0.0.1:%(ENV_VNC_PORT)s
user=headless
autorestart=true
priority=500